# module-contrib/plugsched.service
[Unit]
Description=The plugsched service
ConditionKernelCommandLine=!plugsched.enable=0

[Service]
Type=oneshot
ExecStart=/usr/local/bin/plugsched-install

[Install]
WantedBy=multi-user.target
